# deploy/copilot-dashboard.service
[Unit]
Description=Copilot Usage Dashboard
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/copilot-dashboard
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/copilot-dashboard/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
